# templates/zfs-auto-snapshot.cron.hourly.epp
#!/bin/sh

# Only call zfs-auto-snapshot if it's available
which zfs-auto-snapshot > /dev/null || exit 0

exec zfs-auto-snapshot --quiet --syslog --label=hourly --keep=<%= $zfs_auto_snapshot::keep_hourly_snaps %> //
